# com.jbae.anthropic-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jbae.anthropic-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jbae/Dev/anthropic-digest/anthropic_digest.py</string>
    </array>

    <!-- Run every Monday at 8:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jbae/Dev/anthropic-digest/digest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jbae/Dev/anthropic-digest/digest.log</string>

    <!-- Do NOT run immediately on load -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
